# contrib/bcd-proxy.conf
# Bytesized Proxy
description "Bytesized Connect Proxy"
start on startup

exec /usr/local/bin/bcd-proxy

respawn
respawn limit 5 60
